# database_version/PHP/crontab.cron
*/30 * * * * /usr/bin/php /path/to/your/hosted/hailsDBCompressCron.php >/dev/null 2>&1
